# systemd/pa-shutdown-proxmox.service
[Unit]
Description=Gracefully shutdown all VMs, LXCs, and host
DefaultDependencies=no
Before=shutdown.target halt.target reboot.target
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/pa-shutdown-proxmox.sh --execute
TimeoutStartSec=900
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=halt.target
WantedBy=shutdown.target
WantedBy=reboot.target
